# packaging/monitoring-agent.service

[Unit]
Description=CheckCle Server Monitoring Agent
Documentation=https://github.com/operacle/checkcle-monitoring-agent
After=network.target
Wants=network.target docker.service

[Service]
Type=simple
User=monitoring-agent
Group=monitoring-agent
# Remove docker group requirement - will be handled programmatically if needed
ExecStart=/usr/bin/monitoring-agent
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monitoring-agent

# Environment variables with full PATH
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Security settings (relaxed for Docker access)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/monitoring-agent /var/run/docker.sock
PrivateTmp=true
PrivateDevices=false
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true

# Environment - Load from env file
EnvironmentFile=-/etc/monitoring-agent/monitoring-agent.env
WorkingDirectory=/var/lib/monitoring-agent

[Install]
WantedBy=multi-user.target