# packaging/monitoring-agent.service

[Unit]
Description=CheckCle Server Monitoring Agent
Documentation=https://github.com/checkcle/server-monitoring-agent
After=network.target
Wants=network.target

[Service]
Type=simple
User=monitoring-agent
Group=monitoring-agent
ExecStart=/usr/bin/monitoring-agent
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monitoring-agent

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/monitoring-agent
PrivateTmp=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true

# Environment - Load from env file
EnvironmentFile=-/etc/monitoring-agent/monitoring-agent.env
WorkingDirectory=/var/lib/monitoring-agent

[Install]
WantedBy=multi-user.target