# packaging/postinst

#!/bin/bash
set -e

# Create monitoring-agent user and group if they don't exist
if ! getent group monitoring-agent >/dev/null 2>&1; then
    groupadd --system monitoring-agent
fi

if ! getent passwd monitoring-agent >/dev/null 2>&1; then
    useradd --system --gid monitoring-agent --home-dir /var/lib/monitoring-agent \
            --shell /bin/false --comment "Monitoring Agent" monitoring-agent
fi

# Create necessary directories
mkdir -p /var/lib/monitoring-agent
mkdir -p /var/log/monitoring-agent
mkdir -p /etc/monitoring-agent

# Set ownership and permissions
chown monitoring-agent:monitoring-agent /var/lib/monitoring-agent
chown monitoring-agent:monitoring-agent /var/log/monitoring-agent
chown -R monitoring-agent:monitoring-agent /etc/monitoring-agent

# Set proper permissions
chmod 755 /var/lib/monitoring-agent
chmod 755 /var/log/monitoring-agent
chmod 755 /etc/monitoring-agent
chmod 644 /etc/monitoring-agent/monitoring-agent.env
chmod 644 /etc/monitoring-agent/monitoring-agent.env.example

# Make binary executable
chmod +x /usr/bin/monitoring-agent

# Reload systemd
systemctl daemon-reload

# Display configuration message
echo "======================================================================"
echo "Monitoring Agent installed successfully!"
echo ""
echo "IMPORTANT: Please configure the agent before starting:"
echo "1. Edit /etc/monitoring-agent/monitoring-agent.env"
echo "2. Set required values:"
echo "   - POCKETBASE_URL=your-pocketbase-url"
echo "   - SERVER_TOKEN=your-unique-server-token"
echo "   - SERVER_NAME=your-server-name (optional, will use hostname)"
echo ""
echo "3. Enable and start the service:"
echo "   sudo systemctl enable monitoring-agent"
echo "   sudo systemctl start monitoring-agent"
echo ""
echo "4. Check status:"
echo "   sudo systemctl status monitoring-agent"
echo "   sudo journalctl -u monitoring-agent -f"
echo "======================================================================"