# packaging/postinst

#!/bin/bash
set -e

# Create monitoring-agent user and group
if ! getent group monitoring-agent >/dev/null 2>&1; then
    addgroup --system monitoring-agent
fi

if ! getent passwd monitoring-agent >/dev/null 2>&1; then
    adduser --system --home /var/lib/monitoring-agent --shell /bin/false \
            --gecos "Monitoring Agent" --ingroup monitoring-agent monitoring-agent
fi

# Create directories
mkdir -p /var/lib/monitoring-agent
mkdir -p /var/log/monitoring-agent
mkdir -p /etc/monitoring-agent

# Set ownership and permissions
chown monitoring-agent:monitoring-agent /var/lib/monitoring-agent
chown monitoring-agent:monitoring-agent /var/log/monitoring-agent
chmod 755 /var/lib/monitoring-agent
chmod 755 /var/log/monitoring-agent

# Create environment file if it doesn't exist
if [ ! -f /etc/monitoring-agent/monitoring-agent.env ]; then
    # Check if the config exists in the package location
    if [ -f /etc/monitoring-agent/monitoring-agent.conf ]; then
        cp /etc/monitoring-agent/monitoring-agent.conf /etc/monitoring-agent/monitoring-agent.env
    else
        # Create a basic config file if none exists
        cat > /etc/monitoring-agent/monitoring-agent.env << 'EOF'
# Server Monitoring Agent Configuration
# Auto-generated on installation

# Basic Configuration
AGENT_ID=monitoring-agent-001
CHECK_INTERVAL=30s
HEALTH_CHECK_PORT=9091

# HTTP REST API (fallback)
SERVER_URL=http://localhost:8080
# API_KEY=your-api-key-here

# PocketBase Configuration
POCKETBASE_ENABLED=true
POCKETBASE_URL=http://localhost:8090

# Remote Control
REMOTE_CONTROL_ENABLED=true
COMMAND_CHECK_INTERVAL=10s

# Monitoring Settings
REPORT_INTERVAL=5m
MAX_RETRIES=3
REQUEST_TIMEOUT=10s

# Server Configuration - REQUIRED for proper server registration
SERVER_NAME=
HOSTNAME=
IP_ADDRESS=
OS_TYPE=
SERVER_TOKEN=
EOF
        echo "Created default configuration file"
    fi
fi

# Set permissions for config files
chown root:monitoring-agent /etc/monitoring-agent/monitoring-agent.env
chmod 640 /etc/monitoring-agent/monitoring-agent.env

# Set executable permissions for the binary
if [ -f /usr/bin/monitoring-agent ]; then
    chmod +x /usr/bin/monitoring-agent
    echo "Binary permissions set successfully"
else
    echo "Warning: Binary /usr/bin/monitoring-agent not found"
    ls -la /usr/bin/monitoring-* || echo "No monitoring-* binaries found"
fi

# Reload systemd
systemctl daemon-reload

echo "Server Monitoring Agent installed successfully!"
echo ""
echo "⚠️  IMPORTANT: Server configuration required!"
echo ""
echo "Before starting the service, you must configure the server settings:"
echo "1. Edit /etc/monitoring-agent/monitoring-agent.env"
echo "2. Set the following variables:"
echo "   - SERVER_NAME=your-server-name"
echo "   - HOSTNAME=your-hostname"  
echo "   - IP_ADDRESS=your-server-ip"
echo "   - OS_TYPE=your-os-type"
echo "   - SERVER_TOKEN=your-unique-token"
echo "   - POCKETBASE_URL=your-pocketbase-url"
echo ""
echo "To start the service after configuration:"
echo "   sudo systemctl enable monitoring-agent"
echo "   sudo systemctl start monitoring-agent"
echo ""
echo "To check status:"
echo "   sudo systemctl status monitoring-agent"
echo ""
echo "Health check will be available at: http://localhost:9091/health"